# .crontab
30 * * * * /usr/local/bin/terminal-notifier -title "+= 30" -message ""
