# .crontab
30 * * * * /usr/local/bin/terminal-notifier -title "+= 60" -message ""
0 7 * * * open http://news.ycombinator.com http://twitter.com http://reddit.com http://facebook.com http://pitchfork.com http://nyt.com